# systemd/workspace-hub.service
[Unit]
Description=Workspace Hub Dashboard
After=network.target

[Service]
Type=simple
User=%USER%
WorkingDirectory=%PROJECT_ROOT%/dashboard
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
Environment=NODE_ENV=production
Environment=CONFIG_PATH=%PROJECT_ROOT%/config/projects.json
Environment=DASHBOARD_PORT=8080

[Install]
WantedBy=multi-user.target
